# environments/production/modules/docker/templates/etc/systemd/system/docker.service.d/docker.conf.erb
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon \
<% @docker_opts.each do |a| -%>
    <%= a %> \
<% end -%>
    --host=fd://
